# onecloud/files/etc/crontabs/root
0 * * * * /etc/cn.sh
0 * * * * /etc/gfw.sh
